# docker/supervisord.conf
; ── CatGPT Supervisor Configuration ─────────────────────────────
; Manages Xvfb, x11vnc, noVNC, and the FastAPI server
; All stdout/stderr goes to /dev/fd/1 (Docker logs)

[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

; ── 1. Xvfb — Virtual framebuffer (fake display) ───────────────
[program:xvfb]
command=Xvfb :99 -screen 0 %(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)sx%(ENV_DISPLAY_DEPTH)s -ac -nolisten tcp +extension GLX
autorestart=true
priority=10
stdout_logfile=/app/logs/xvfb.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/app/logs/xvfb.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
startsecs=2

; ── 2. x11vnc — VNC server capturing Xvfb display ──────────────
[program:x11vnc]
command=x11vnc -display :99 -forever -shared -rfbport 5900 -xkb -noncache -rfbauth /app/.vnc/passwd
autorestart=true
priority=20
startretries=10
startsecs=3
stdout_logfile=/app/logs/x11vnc.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/app/logs/x11vnc.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

; ── 3. noVNC — Browser-based VNC viewer ─────────────────────────
[program:novnc]
command=websockify --web=/usr/share/novnc/ 6080 localhost:5900
autorestart=true
priority=30
startsecs=3
stdout_logfile=/app/logs/novnc.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/app/logs/novnc.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

; ── 4. FastAPI Server — CatGPT API ─────────────────────────────
[program:catgpt]
command=python -m src.api.server
directory=/app
autorestart=true
priority=40
startsecs=10
startretries=3
stopwaitsecs=30
stdout_logfile=/app/logs/catgpt_docker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/catgpt_docker.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=DISPLAY=":99",PYTHONUNBUFFERED="1",PYTHONPATH="/app"
